# macosx/crontab/beefy.local
# Backup Jobs
# min hour dom mon dow cmd
#  20  03   *   *    *  sh $HOME/Develop/pmade/rc/backup/slim.pmade.com/postgresql.sh
#  30  04   *   *    *  sh $HOME/Develop/pmade/rc/backup/slim.pmade.com/subversion.sh

# Morning Wake-Up Schedule
# min hour dom mon dow cmd
   00  06   *   *    *  sh $HOME/Develop/pmade/rc/macosx/scripts/iTunes/kgnu-rip.sh
   30  06   *   *    *  osascript $HOME/Develop/pmade/rc/macosx/scripts/iTunes/wake-up.applescript
   35  06   *   *    *  sh $HOME/Develop/pmade/rc/macosx/scripts/iTunes/morning-news.sh
   45  06   *   *    *  sh $HOME/Develop/pmade/rc/macosx/scripts/iTunes/morning-energy.sh
